# compile/files/etc/rc.local
# Put your custom commands here that should be executed once
# x-gpio only for hg680-p & b860h
#sleep 7 && /usr/bin/x-gpio -r
#sleep 7 && /sbin/free.sh
#/sbin/jam bug.com
exit 0